# systemd/tor-manager.service
[Unit]
Description=Tor Manager - Tor Proxy Service
Documentation=https://www.torproject.org/
After=network-online.target nss-lookup.target
Wants=network-online.target

[Service]
Type=notify
User=root
ExecStart=/opt/tor/bin/tor -f /etc/tor/torrc
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30

# 安全加固
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/
ReadWritePaths=/var/lib/tor /var/log/tor /root/self-check-tor-systemd/var/run
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID

# 资源限制
LimitNOFILE=65535
LimitNPROC=512

# 环境变量
Environment=TOR_SOCKS_PORT=9050
Environment=TOR_CONTROL_PORT=9051

# PID 文件
PIDFile=/root/self-check-tor-systemd/var/run/tor.pid

[Install]
WantedBy=multi-user.target
